# deploy/systemd/chatko-scheduler.service
[Unit]
Description=Chatko Laravel Scheduler Worker
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/chatko/backend
ExecStart=/usr/bin/php artisan schedule:work
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
